# ova/rootfs/etc/init.d/wlcsim
#!/sbin/openrc-run

name="wlcsim"
description="Cisco 9800-CL WLC Simulator"
command="/usr/local/bin/wlcsim"
command_args="--config /etc/wlcsim/devices.yaml --lan --dashboard-port 8080"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/wlcsim.log"
error_log="/var/log/wlcsim.log"

depend() {
    after firewall local
}

start_pre() {
    checkpath -d -m 0755 /var/log
    # Wait for DHCP to install a default route (up to 30s). The dhcp-all.start
    # local script runs udhcpc in the background, so routes may not exist yet
    # when this service fires.
    i=0
    while [ $i -lt 30 ]; do
        if ip route show default 2>/dev/null | grep -q '^default'; then
            return 0
        fi
        sleep 1
        i=$((i + 1))
    done
    ewarn "No default route after 30s; starting anyway"
}
